# supervisor-app.conf
[supervisord]
nodaemon=true

[program:web-ui]
directory = /app/
command = npm start
autostart = true
autorestart = false
