# web-ui/supervisor-app.conf
[supervisord]
nodaemon=true

[program:web-ui]
command = npm start;
autostart=true
autorestart=false
startsecs = 0
redirect_stderr=true
